# docker/hdfs/rootfs/etc/s6-overlay/s6-rc.d/kadmin/finish
#!/usr/bin/with-contenv bash
set -eo pipefail
echo "[INFO] $(basename $(dirname $0)) service stopping (exit: ${1}, signal: ${2})"
exit 0
