# docker/hdfs/rootfs/etc/s6-overlay/s6-rc.d/kadmin/run
#!/usr/bin/with-contenv bash
exec /usr/sbin/kadmind -nofork
