# docker/hdfs/rootfs/etc/s6-overlay/s6-rc.d/datanode/run
#!/usr/bin/with-contenv bash
export JAVA_HOME=/usr/lib/jvm/java-21-openjdk-amd64
export HADOOP_HOME=/opt/hadoop
export HADOOP_OPTS="-Djava.security.krb5.conf=/etc/krb5.conf"

# Obtain Kerberos ticket using the container's hostname (set via
# HDFS_HOSTNAME build arg). Falls back to localhost for pre-existing
# images that don't set the env var.
: "${HDFS_HOSTNAME:=localhost}"
gosu mntrs /usr/bin/kinit -kt /tmp/hdfs.keytab "hdfs/${HDFS_HOSTNAME}@TEST.LOCAL" 2>/dev/null || true

# Wait for NameNode. Use the configured defaultFS (hdfs://nameservice
# in HA mode, or whatever core-site.xml specifies) rather than
# hardcoding localhost:8020 – the latter bypasses the HA proxy and
# can fail intermittently when the NameNode is registered under the
# nameservice rather than as a standalone host.
echo "[datanode] Waiting for NameNode..."
for i in $(seq 1 30); do
    if gosu mntrs /opt/hadoop/bin/hdfs dfsadmin -report >/dev/null 2>&1; then
        echo "[datanode] NameNode is ready"
        break
    fi
    sleep 2
done

echo "[datanode] Starting..."
exec gosu mntrs /opt/hadoop/bin/hdfs datanode
