# docker/hdfs/rootfs/etc/s6-overlay/s6-rc.d/datanode/run
#!/usr/bin/with-contenv bash
export JAVA_HOME=/usr/lib/jvm/java-21-openjdk-amd64
export HADOOP_HOME=/opt/hadoop
export HADOOP_OPTS="-Djava.security.krb5.conf=/etc/krb5.conf"

# Obtain Kerberos ticket using the container's hostname (set via
# HDFS_HOSTNAME build arg). Falls back to localhost for pre-existing
# images that don't set the env var.
: "${HDFS_HOSTNAME:=localhost}"
gosu mntrs /usr/bin/kinit -kt /tmp/hdfs.keytab "hdfs/${HDFS_HOSTNAME}@TEST.LOCAL" 2>/dev/null || true

# Wait for NameNode. Use -fs hdfs://localhost:8020 to bypass the HA
# proxy — at container boot the proxy may race s6-overlay startup.
# The Datanode binary itself uses the configured defaultFS (hdfs-site.xml)
# for its actual operation, so this wait path doesn't affect how it
# connects once running.
echo "[datanode] Waiting for NameNode..."
for i in $(seq 1 30); do
    if gosu mntrs /opt/hadoop/bin/hdfs dfsadmin -fs hdfs://localhost:8020 -report >/dev/null 2>&1; then
        echo "[datanode] NameNode is ready"
        break
    fi
    sleep 2
done

echo "[datanode] Starting..."
exec gosu mntrs /opt/hadoop/bin/hdfs datanode
